# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/BLAOS/system/BLAOS/current/BLAOS"
command_args="config"