# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/haos/system/haos/current/haos"
command_args="config"